# scripts/com.cloudcodetree.feed-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Auto-publish the AI News blog when the syndication feed changes.

  The Claude Desktop "AI Developer News" task writes content/feed.xml into the
  connected repo folder but cannot push. This LaunchAgent watches that file and,
  on change, runs scripts/push-feed.sh (ingest → commit → push → deploy).

  WatchPaths fires within seconds of the task writing the feed. StartInterval is
  an hourly safety net in case a file event is missed; push-feed.sh is idempotent,
  so an extra run with no feed change commits nothing.

  Install:
    cp scripts/com.cloudcodetree.feed-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.cloudcodetree.feed-sync.plist
  Disable:
    launchctl unload ~/Library/LaunchAgents/com.cloudcodetree.feed-sync.plist
  Logs: /tmp/cct-feed-sync.log, /tmp/cct-feed-sync.err
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cloudcodetree.feed-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/chris.harper/Development/cloudcodetree.github.io/scripts/push-feed.sh</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/chris.harper/Development/cloudcodetree.github.io/content/feed.xml</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/tmp/cct-feed-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cct-feed-sync.err</string>
</dict>
</plist>
